# api/ratewatch.service
[Unit]
Description=Ohio Rate Watch API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ratewatch-api
EnvironmentFile=/opt/ratewatch-api/.env
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
